# fan-curve/msi-fan-curve.service
[Unit]
Description=MSI fan curve (msi-wmi-platform)
# No ordering on the module: the driver probes asynchronously and the apply
# script waits for the hwmon device itself. Resume needs no hook either -- the
# driver re-applies the fan state after firmware resume (deep S3/hibernate).

[Service]
Type=oneshot
ExecStart=/usr/local/bin/msi-fan-curve /etc/msi-fan-curve.conf
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
